# nzbget.ini
[program:NZBGet]
user=
command=/usr/bin/nzbget -c /config/nzbget.conf -s
